# etc/services.d/redis/run
#!/usr/bin/with-contenv sh

set -eo pipefail

cd / || exit 1

/usr/bin/redis-server \
  > /dev/null 2>&1
  #2>&1 | awk '{print "[redis] " $0}'
  #> /dev/null 2>&1
  #

